# systemd/fbcp.service
[Unit]
Description=Framebuffer mirror (HDMI → TFT via SPI)
# Start early — after local-fs and udev so /dev/spidev0.0 exists.
# Before the display manager so boot messages are visible on TFT.
DefaultDependencies=no
After=local-fs.target systemd-udev-settle.service
Before=display-manager.service graphical.target
Wants=local-fs.target

[Service]
Type=simple
# Wait for SPI device to appear (up to 10 seconds)
ExecStartPre=/bin/bash -c '\
  for i in $(seq 1 10); do \
    [ -e /dev/spidev0.0 ] && exit 0; \
    sleep 1; \
  done; \
  modprobe spi_bcm2835 2>/dev/null || true; \
  modprobe spidev 2>/dev/null || true; \
  sleep 2; \
  [ -e /dev/spidev0.0 ]'
# Wait for /dev/fb0 (vc4 framebuffer may take a moment)
ExecStartPre=/bin/bash -c '\
  for i in $(seq 1 15); do \
    [ -e /dev/fb0 ] && exit 0; \
    sleep 1; \
  done; \
  echo "Warning: /dev/fb0 not found after 15s" >&2'
ExecStart=/usr/local/bin/fbcp --fps=10 --touch --touch-swap-xy --touch-invert-x
Restart=on-failure
RestartSec=5
KillMode=process
TimeoutStopSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
