# systemd/ili9481-fb.service
[Unit]
Description=ILI9481 userspace framebuffer driver
After=local-fs.target systemd-modules-load.service
# Safety: never reboot/halt/etc on repeated failure — just stop retrying.
StartLimitBurst=5
StartLimitIntervalSec=120
StartLimitAction=none

[Service]
Type=simple
# Force HDMI connector on so KMS creates /dev/fb0 (even without a
# physical HDMI monitor).  The kernel cmdline 'video=HDMI-A-1:…D' does
# this too, but writing to sysfs is a belt-and-suspenders fallback.
ExecStartPre=/bin/bash -c '\
  for i in 1 2 3 4 5; do \
    if [ -e /sys/class/drm/card0-HDMI-A-1/status ]; then \
      echo on > /sys/class/drm/card0-HDMI-A-1/status 2>/dev/null || true; \
      break; \
    fi; \
    sleep 1; \
  done'
# Wait for /dev/fb0 to appear after connector is forced on
ExecStartPre=/bin/bash -c '\
  for i in 1 2 3 4 5 6 7 8 9 10; do \
    [ -e /dev/fb0 ] && exit 0; \
    sleep 1; \
  done; \
  echo "Warning: /dev/fb0 not found after 10s" >&2'
ExecStart=/usr/local/bin/ili9481-fb --config=/etc/ili9481/ili9481.conf
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
