# systemd/ili9481-fb.service
[Unit]
Description=ILI9481 userspace framebuffer driver
After=local-fs.target
# Safety: never reboot/halt/etc on repeated failure — just stop retrying.
StartLimitBurst=3
StartLimitIntervalSec=120
StartLimitAction=none

[Service]
Type=simple
# Give the kernel time to initialise /dev/fb0 and /dev/gpiomem before we grab them.
ExecStartPre=/bin/sleep 3
ExecStart=/usr/local/bin/ili9481-fb --config=/etc/ili9481/ili9481.conf
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
